# deploy/openwcs-deploy.timer
[Unit]
Description=Poll origin/main and redeploy openWCS when it changes

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
Unit=openwcs-deploy.service

[Install]
WantedBy=timers.target
